# config/mac/Library/LaunchDaemons/io.github.gershnik.wsddn.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.github.gershnik.wsddn</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/wsddn</string>
        <string>--launchd</string>
        <string>--config=/etc/wsddn.conf</string>
        <string>--log-os-log</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/var/empty</string>
    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
    </dict>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
