# config/mac/Library/LaunchDaemons/io.github.gershnik.wsddn.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.github.gershnik.wsddn</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/wsddn</string>
        <string>--launchd</string>
        <string>--config=/etc/wsddn.conf</string>
        <string>--pid-file=/var/run/wsddn/wsddn.pid</string>
        <string>--log-file=/var/log/wsddn.log</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/var/empty</string>
    <key>EnvironmentVariables</key>
	<dict>
	    <key>PATH</key>
	    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/usr/local/sbin</string>
	</dict>
    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/log/wsddn.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/wsddn.log</string>
</dict>
</plist>
